# supervisord.conf
[supervisord]
nodaemon = true
logfile = /var/log/supervisor/supervisord.log
childlogdir = /var/log/supervisor
logfile_maxbytes = 10mb

[include]
files = /etc/supervisor/conf.d/*.conf

[program::sshd]
command = /usr/sbin/sshd -D

[program::mysql]
command=mysqld_safe
autorestart=true